# deploy/systemd/api.service.in
[Unit]
Description=Wabbit API Server
Requires=mysqld.service
After=mysqld.service

[Service]
Type=simple
User=@USER@
Group=@GROUP@
WorkingDirectory=@WORKING_DIR@
ExecStart=@GUNICORN_BIN@ @APP_MODULE@ --workers @gunicorn_workers@ --threads @gunicorn_threads@ --bind 127.0.0.1:@gunicorn_port@
Restart=always

[Install]
WantedBy=multi-user.target
