# systemd/wabbit-api.service.in
[Unit]
Description=Wabbit API Server
Requires=mysqld.service
After=mysqld.service

[Service]
Type=simple
User=@USER@
Group=@GROUP@
WorkingDirectory=@WORKING_DIR@
ExecStart=@GUNICORN_BIN@ @APP_MODULE@ -b @ADDRESS@

[Install]
WantedBy=multi-user.target
