# systemd/wabbit-record.service.in
[Unit]
Description=Wabbit Recorder
Requires=mysqld.service
After=mysqld.service

[Service]
Type=simple
User=@USER@
Group=@GROUP@
WorkingDirectory=@WORKING_DIR@
ExecStart=@RECORD_CMD@ -1
Restart=always

[Install]
WantedBy=multi-user.target
